# roles/forge_backup/templates/forge-restic.service.j2
[Unit]
Description=forge restic off-site backup
Wants=network-online.target
After=network-online.target zfs.target

[Service]
Type=oneshot
EnvironmentFile=/etc/forge/restic/env
ExecStart=/usr/local/sbin/forge-restic-backup
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
