# roles/forge_backup/templates/syncoid.timer.j2
[Unit]
Description=forge syncoid replication schedule: {{ item.1.source }}

[Timer]
OnCalendar={{ item.1.schedule | default('*-*-* *:15:00') }}
Persistent=true

[Install]
WantedBy=timers.target
